# deploy/systemd/dnd-backend.service
[Unit]
Description=D&D Knowledge Base FastAPI Backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/dnd-knowledge-base/backend
Environment=PATH=/opt/dnd-knowledge-base/venv/bin
EnvironmentFile=/opt/dnd-knowledge-base/.env
ExecStart=/opt/dnd-knowledge-base/venv/bin/uvicorn backend.main:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 2 \
    --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
